# templates/systemd.service
[Unit]
Description=Prometheus exporter : {{EXPORTER_NAME}}
After=network.target

[Service]
Type=simple
User={{EXPORTER_NAME}}
Group={{EXPORTER_NAME}}
Restart=always
EnvironmentFile=/etc/default/{{EXPORTER_NAME}}
ExecStart=/usr/bin/{{EXPORTER_NAME}} $ARGS

[Install]
WantedBy=multi-user.target
